# check_mk_agent.socket
[Unit]
Description=Check_MK Agent Socket

[Socket]
ListenStream=6556
Accept=yes

[Install]
WantedBy=sockets.target

